# iptables-nakivo/etc/rc.local
#!/bin/bash

# script para inicializar regras iptables para o Nakivo
/usr/local/bin/iptables-nakivo.sh

exit 0
