# etc/rc.local
#!/bin/bash

/usr/local/bin/iptables-nakivo.sh

exit 0
